# platform/roles/herolib_backend/templates/herolib.service.j2
[Unit]
Description=Herolib Backend Server
After=network.target redis-6379.service redis-6381.service

[Service]
Type=simple
User=hero
Group=hero
WorkingDirectory=/opt/herolib
ExecStart=/usr/local/bin/v run /opt/herolib/examples/hero/heromodels/heroserver_example.vsh
Restart=always
Environment=HOME=/opt/hero
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=DOMAIN_NAME={{ domain_name | default('') }}
Environment=ENABLE_SSL={{ enable_ssl | default('false') }}
Environment=REDIS_PASSWORD={{ redis_password | default('') }}

[Install]
WantedBy=multi-user.target